# templates/default/docker-cgroup-rewrite.service.erb
[Unit]
Description=docker-cgroup-rewrite
After=containerd.service

[Service]
ExecStart=<%= @docker_cgroup_rewrite_script %> \
memory_limit_bytes=<%= @memory_limit_bytes %> \
memory_soft_limit_bytes=<%= @memory_soft_limit_bytes %> \
cpu_quota=<%= @cpu_quota %>

Restart=always
Delegate=yes
# Having non-zero Limit*s causes performance problems due to accounting overhead
# in the kernel. We recommend using cgroups to do container-local accounting.
LimitNPROC=infinity
LimitCORE=infinity
LimitNOFILE=infinity
# Comment TasksMax if your systemd version does not supports it.
# Only systemd 226 and above support this version.
TasksMax=infinity

[Install]
WantedBy=multi-user.target
