# templates/default/namenode.service.erb
[Unit]
Description = NameNode server for HDFS.
After = syslog.target network.target remote-fs.target <%= @deps %>
#The namenode takes ~32s to fail if glassfish is not up
#restartSec is 2 so the service restart every 34s
#a maximum burst of 18 in 11min guaranty that it will at least try for 10min (18*34=612s=10,2min)
StartLimitIntervalSec=660
<% if node['hops']['alive-watchdog']['enabled'].casecmp?("true") -%>
# Do not be very strict as it is normal for the service to fail and restart immediatelly
# when liveness watchdog is enabled
StartLimitBurst=350
<% else -%>
StartLimitBurst=18
<% end -%>

[Service]
User = <%= node['hops']['hdfs']['user'] %>

PIDFile = <%= node['hops']['logs_dir'] %>/hadoop-<%= node['hops']['hdfs']['user'] %>-namenode.pid
ExecStart = <%= node['hops']['base_dir'] %>/sbin/start-nn.sh
ExecStop = <%= node['hops']['base_dir'] %>/sbin/stop-nn.sh
ExecReload = <%= node['hops']['base_dir'] %>/sbin/restart-nn.sh
LimitNOFILE = 65000
LimitNPROC = 10000
Restart=on-failure
RestartSec=2

Type = forking

[Install]
WantedBy = multi-user.target
