# templates/default/namenode.service.erb
[Unit]
Description = NameNode server for HDFS.
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
User = <%= node['hops']['hdfs']['user'] %>

PIDFile = <%= node['hops']['logs_dir'] %>/hadoop-<%= node['hops']['hdfs']['user'] %>-namenode.pid
ExecStart = <%= node['hops']['base_dir'] %>/sbin/start-nn.sh
ExecStop = <%= node['hops']['base_dir'] %>/sbin/stop-nn.sh
ExecReload = <%= node['hops']['base_dir'] %>/sbin/restart-nn.sh
LimitNOFILE = 65000
LimitNPROC = 10000
Restart=on-failure
RestartSec=2

Type = forking

[Install]
WantedBy = multi-user.target
