# templates/default/nodemanager.service.erb
[Unit]
Description = NodeManager. The Processing Nodes for YARN.
After = syslog.target network.target remote-fs.target <%= @deps %>
#The nodemanager takes ~36s to fail if the resourcemanager is not up
#restartSec is 5 so the service restart every 41s
#a maximum burst of 15 in 11min guaranty that it will at least try for 10min (15*41=615s=10,25min)
StartLimitIntervalSec=660
StartLimitBurst=15

[Service]
User = <%= node['hops']['yarn']['user'] %>
Group = <%= node['hops']['group'] %>

PIDFile = <%= node['hops']['logs_dir'] %>/hadoop-<%= node['hops']['yarn']['user'] %>-nodemanager.pid

ExecStart = <%= node['hops']['base_dir'] %>/sbin/start-nm.sh
ExecStop = <%= node['hops']['base_dir'] %>/sbin/stop-nm.sh
ExecReload = <%= node['hops']['base_dir'] %>/sbin/restart-nm.sh
LimitNOFILE = 65000
LimitNPROC = 10000
TimeoutSec= 10min
Restart=on-failure
RestartSec=5

Type = forking

[Install]
WantedBy = multi-user.target
