# templates/default/resourcemanager.service.erb
[Unit]
Description = ResourceManager. The scheduler for YARN.
After = syslog.target network.target remote-fs.target <%= @deps %>

[Service]
User = <%= node['hops']['rm']['user'] %>

PIDFile = <%= node['hops']['logs_dir'] %>/hadoop-<%= node['hops']['rm']['user'] %>-resourcemanager.pid
ExecStart = <%= node['hops']['base_dir'] %>/sbin/start-rm.sh
ExecStop = <%= node['hops']['base_dir'] %>/sbin/stop-rm.sh
ExecReload = <%= node['hops']['base_dir'] %>/sbin/restart-rm.sh
LimitNOFILE = 65000
LimitNPROC = 10000

Type = forking

[Install]
WantedBy = multi-user.target
